# deploy/tbox_server.service
[Unit]
Description=Home Updater Service
Documentation=https://github.com/your-repo/tbox
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/usr/local/tbox
ExecStart=/usr/local/tbox/bin/tbox_server
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/usr/local/tbox/logs /usr/local/tbox/data

# Environment
Environment="TBOX_LOG_DIR=/usr/local/tbox/logs"

# Resource limits
LimitNOFILE=65536
LimitCORE=infinity

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tbox

[Install]
WantedBy=multi-user.target
